# .config/systemd/user/next-bg.service
[Unit]
Description=Changes background
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/next-bg

[Install]
WantedBy=default.target
